# server-image/supervisord.conf
[supervisord]
nodaemon=true

[program:hma]
environment=SERVER=_SERVER_
command=/usr/bin/hma-runner.sh

[program:privoxy]
command=/usr/sbin/privoxy --no-daemon /etc/privoxy/config
